# HW4/crontab
# run nginx log analyser every hour
0 1 * * * 	./nginx_la.sh
