# HW4/crontab
# run nginx log analyser eveery hour
0 1 * * * 	./nginx_la.sh
